# config/server/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:apache2]
priority=1
command=httpd -D FOREGROUND
